# pi-hole/rootfs/etc/services.d/dnscrypt-proxy/run
#!/bin/bash

set -ueo pipefail

/dnscrypt-proxy
